# Morning_report/com.yuichiromori.morningreport.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yuichiromori.morningreport</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/yuichiromori/Desktop/NCAフォルダ/morning_report/run.sh</string>
    </array>

    <!-- 毎朝 7:00 に実行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/yuichiromori/Desktop/NCAフォルダ/morning_report/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yuichiromori/Desktop/NCAフォルダ/morning_report/launchd_stderr.log</string>

    <!-- スリープ中に起動時刻を過ぎた場合、起動後に実行 -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
